# sentinix/usr/src/rc_scripts/rc.d/rc.local
#!/bin/sh
#
# /etc/rc.d/rc.local:  Local system initialization script.
#
# Put any local setup commands in here:

/usr/bin/kbdrate -r 30 -d 250

