# sentinix/usr/src/syswideconfig/build.crontab
#!/bin/sh

cd misc &&

cat run-parts > /usr/bin/run-parts &&
chmod 0755 /usr/bin/run-parts &&
crontab - -u root < root.crontab
