# sentinix/usr/src/syswideconfig/build.crontab
#!/bin/sh

zcat run-parts.gz > /usr/bin/run-parts &&
chmod 0755 /usr/bin/run-parts &&
crontab - -u root < root.crontab
